# packaging/loguard-check.timer
[Unit]
Description=Loguard Watchdog Timer
Requires=loguard.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target